# deploy/vcpulse.service
[Unit]
Description=VC Pulse Dashboard
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/vcpulse
ExecStart=/home/ubuntu/vcpulse/venv/bin/uvicorn app:app --host 127.0.0.1 --port 8003
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
